# cachyos-rate-mirrors/cachyos-rate-mirrors.timer
[Unit]
Description=Automatically refresh arch / cachyos mirror lists

[Timer]
OnBootSec=15s
OnUnitActiveSec=2d
Persistent=true

[Install]
WantedBy=timers.target
