# cachyos-rate-mirrors/cachyos-rate-mirrors.timer
[Unit]
Description=Automatically refresh arch / cachyos mirror lists

[Timer]
OnBootSec=2min
RandomizedDelaySec=30s
OnUnitActiveSec=5d
Persistent=true

[Install]
WantedBy=timers.target
