# cachyos-rate-mirrors/cachyos-rate-mirrors.timer
[Unit]
Description=Automatically refresh arch / cachyos mirror lists
Wants=network-online.target
After=network-online.target

[Timer]
OnUnitActiveSec=10d
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
