# cachyos-rate-mirrors/cachyos-rate-mirrors.timer
[Unit]
Description=Weekly refresh arch / cachyos mirror lists

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target
