# valkey/valkey-sentinel.service
[Unit]
Description=High availability daemon for Redis
After=network.target

[Service]
Type=notify
User=valkey
Group=valkey
ExecStart=/usr/bin/valkey-sentinel /etc/valkey/sentinel.conf --supervised systemd
TimeoutStartSec=60
TimeoutStopSec=60
CapabilityBoundingSet=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/valkey/
NoNewPrivileges=true
RuntimeDirectory=valkey
RuntimeDirectoryMode=755
LimitNOFILE=10032

[Install]
WantedBy=multi-user.target
Alias=redis-sentinel.service
